# 10-monitoring/01-base/monitoring.cron
* * * * * root test -f /opt/monitoring/monitoring.py && cd /opt/monitoring && /usr/bin/python3 monitoring.py